# web/docker/8.3/supervisor-queue.conf.example
# Supervisor configuration for Laravel Queue Worker
# This configuration can be copied to /etc/supervisor/conf.d/ for production use
#
# Usage:
#   1. Copy this file: sudo cp supervisor-queue.conf.example /etc/supervisor/conf.d/queue-worker.conf
#   2. Edit the file to set your environment-specific values
#   3. Update supervisor: sudo supervisorctl reread && sudo supervisorctl update
#   4. Start the service: sudo supervisorctl start queue-worker
#
# Configuration:
#   - Replace /var/www/html with your actual Laravel application path
#   - Adjust queues (conversion,default) if you have different queue names
#   - Adjust user to match your web server user (www-data, nginx, etc.)

[program:queue-worker]
command=/usr/bin/php -d variables_order=EGPCS /var/www/html/artisan queue:work --queue=conversion,default
user=www-data
environment=LARAVEL_SAIL="1"
stdout_logfile=/var/log/supervisor/queue-worker.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/queue-worker_err.log
stderr_logfile_maxbytes=10MB
autorestart=true
stopwaitsecs=30
